# ansibles/roles/backup.postgres/templates/postgres-backup.timer.j2
[Unit]
Description=Dump and back up postgres databases

[Timer]
OnCalendar=*-*-* 3:00
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
