# ansibles/roles/backup.cephfs/templates/cephfs-backup.service.j2
[Unit]
Description=Back up CephFS filesystems
After=network.target network-online.target

[Service]
Type=oneshot
ExecStart=/opt/backup/cephfs-backup.sh
